# etc/cron.hourly/zfs-auto-backup-daily-local
#!/bin/sh

# Only call zfs-auto-backup if it's available
which zfs-auto-backup-daily-local > /dev/null || exit 0

exec zfs-auto-backup-daily-local

