# etc/cron.hourly/zfs-auto-backup-daily-remote
#!/bin/sh

# Only call zfs-auto-backup-daily-remote if it's available
which zfs-auto-backup-daily-remote > /dev/null || exit 0

exec zfs-auto-backup-daily-remote

